# supervisor/piui-supervisor
### BEGIN INIT INFO
# Provides: PiUi Supervisor
# Required-Start: $remote_fs $syslog
# Required-Stop: $remote_fs $syslog
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: PiUi Supervisor
# Description: PiUi Supervisor
### END INIT INFO
 
 
#! /bin/sh
# /etc/init.d/piui-supervisor
 
 
export HOME
case "$1" in
    start)
        echo "Starting PiUi Supervisor"
        /home/pi/piui/supervisor/piui-supervisor.py 2>&1 &
    ;;
    stop)
        echo "Stopping PiUi Supervisor"
    LCD_PID=`ps auxwww | grep piui-supervisor.py | head -1 | awk '{print $2}'`
    kill -9 $LCD_PID
    ;;
    *)
        echo "Usage: /etc/init.d/piui-supervisor {start|stop}"
        exit 1
    ;;
esac
exit 0